# config/systemd/openclaw-mem-sync.timer
[Unit]
Description=Run OpenClaw Memory Sync daily at 03:00

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
